# systemd/user/vellumd.service
[Unit]
Description=Vellum Wayland Wallpaper Daemon
Documentation=https://github.com/yourname/vellum
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/vellumd
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info
NoNewPrivileges=true

[Install]
WantedBy=default.target
