# packaging/vellum-daemon.service
[Unit]
Description=Vellum wallpaper session daemon
Documentation=man:systemd.user(5)
PartOf=graphical-session.target
After=graphical-session.target
BindsTo=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=notify
ExecStart=/usr/bin/vellum --daemon --namespace vellum
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
